# systemd/scrolllock-led-daemon.service
[Unit]
Description=Scroll Lock LED Daemon
After=multi-user.target

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/local/bin/scrolllock-led-daemon
Restart=on-failure
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target