# debian/scrolllock-led-daemon.service
[Unit]
Description=Daemon que sincroniza LED de Scroll Lock com teclado
After=network.target

[Service]
Type=notify
ExecStart=/usr/local/bin/scrolllock-led-daemon
# HEARTBEAT CONFIG - CRÍTICO PARA DETECÇÃO DE TRAVAMENTO
WatchdogSec=30
# Reinicia automaticamente em caso de falha ou travamento
Restart=on-failure
RestartSec=5
# Limita recursos para segurança
PrivateTmp=yes
ProtectSystem=strict
NoNewPrivileges=yes
CapabilityBoundingSet=CAP_SYS_TTY_CONFIG CAP_SYS_TTY_READ

[Install]
WantedBy=multi-user.target